# var/lib/dpkg/info/bluez.list
/.
/usr
/usr/lib
/usr/lib/bluetooth
/usr/lib/bluetooth/bluetoothd
/usr/bin
/usr/bin/rfcomm
/usr/bin/hcitool
/usr/bin/rctest
/usr/bin/btmon
/usr/bin/l2ping
/usr/bin/hciattach
/usr/bin/bluetoothctl
/usr/bin/ciptool
/usr/bin/gatttool
/usr/bin/bccmd
/usr/bin/sdptool
/usr/bin/l2test
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.bluez.service
/usr/share/doc
/usr/share/doc/bluez
/usr/share/doc/bluez/README.Debian.gz
/usr/share/doc/bluez/NEWS.Debian.gz
/usr/share/doc/bluez/changelog.gz
/usr/share/doc/bluez/changelog.Debian.armhf.gz
/usr/share/doc/bluez/changelog.Debian.gz
/usr/share/doc/bluez/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/bluetoothd.8.gz
/usr/share/man/man1
/usr/share/man/man1/l2ping.1.gz
/usr/share/man/man1/btmon.1.gz
/usr/share/man/man1/rfcomm.1.gz
/usr/share/man/man1/hcitool.1.gz
/usr/share/man/man1/hciconfig.1.gz
/usr/share/man/man1/sdptool.1.gz
/usr/share/man/man1/hid2hci.1.gz
/usr/share/man/man1/rctest.1.gz
/usr/share/man/man1/hciattach.1.gz
/usr/share/man/man1/l2test.1.gz
/usr/share/man/man1/bccmd.1.gz
/usr/share/man/man1/bluetoothctl.1.gz
/usr/share/man/man1/ciptool.1.gz
/usr/sbin
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bluetooth.service
/lib/udev
/lib/udev/hid2hci
/lib/udev/rules.d
/lib/udev/rules.d/97-hid2hci.rules
/lib/udev/rules.d/50-bluetooth-hci-auto-poweron.rules
/bin
/bin/hciconfig
/etc
/etc/init.d
/etc/init.d/bluetooth
/etc/bluetooth
/etc/bluetooth/proximity.conf
/etc/bluetooth/network.conf
/etc/bluetooth/main.conf
/etc/bluetooth/input.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/bluetooth.conf
/etc/default
/etc/default/bluetooth
/usr/sbin/bluetoothd
