# var/lib/dpkg/info/hddtemp.list
/.
/etc
/etc/init.d
/etc/init.d/hddtemp
/etc/default
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/hddtemp
/etc/hddtemp.db
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/hddtemp.8.gz
/usr/share/locale
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/hddtemp.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/hddtemp.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/hddtemp.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/hddtemp.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/hddtemp.mo
/usr/share/doc
/usr/share/doc/hddtemp
/usr/share/doc/hddtemp/changelog.Debian.gz
/usr/share/doc/hddtemp/changelog.gz
/usr/share/doc/hddtemp/TODO
/usr/share/doc/hddtemp/README
/usr/share/doc/hddtemp/README.Debian
/usr/share/doc/hddtemp/contribs
/usr/share/doc/hddtemp/contribs/wminfo
/usr/share/doc/hddtemp/contribs/README
/usr/share/doc/hddtemp/contribs/crash.c
/usr/share/doc/hddtemp/contribs/hddtemp-all.sh
/usr/share/doc/hddtemp/contribs/analyze
/usr/share/doc/hddtemp/contribs/analyze/graph-field.sh
/usr/share/doc/hddtemp/contribs/analyze/plot-field.pl
/usr/share/doc/hddtemp/contribs/analyze/hddtemp_monitor.sh
/usr/share/doc/hddtemp/copyright
/usr/sbin
/usr/sbin/hddtemp
