# var/lib/dpkg/info/hdparm.conffiles
/etc/apm/event.d/20hdparm
/etc/default/hdparm
/etc/hdparm.conf
/etc/init.d/hdparm
