# var/lib/dpkg/info/lirc.list
/.
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/85-lirc.rules
/etc
/etc/lirc
/etc/lirc/lircmd.conf
/etc/lirc/hardware.conf
/etc/lirc/lircd.conf
/etc/init.d
/etc/init.d/lirc
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/lirc
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/lirc
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.workstation/lirc
/usr
/usr/share
/usr/share/doc-base
/usr/share/doc-base/lirc
/usr/share/lirc
/usr/share/lirc/remotes
/usr/share/lirc/remotes/remotec
/usr/share/lirc/remotes/remotec/lircd.conf.BW6130
/usr/share/lirc/remotes/remotec/lircmd.conf.BW6130
/usr/share/lirc/remotes/pixelview
/usr/share/lirc/remotes/pixelview/lircd.conf.remotemaster
/usr/share/lirc/remotes/pixelview/lircd.conf.playtv_bt878
/usr/share/lirc/remotes/pixelview/lircd.conf.playtv
/usr/share/lirc/remotes/pixelview/lircmd.conf.playtv_pro
/usr/share/lirc/remotes/pixelview/lircd.conf.playtv_pro
/usr/share/lirc/remotes/pixelview/lircmd.conf.remotemaster
/usr/share/lirc/remotes/pixelview/lircrc.playtv_pro
/usr/share/lirc/remotes/life-view
/usr/share/lirc/remotes/life-view/lircrc.flyvideo
/usr/share/lirc/remotes/life-view/lircmd.conf.flyvideo
/usr/share/lirc/remotes/life-view/lircd.conf.flyvideo
/usr/share/lirc/remotes/x10
/usr/share/lirc/remotes/x10/lircmd.conf.mouseremote
/usr/share/lirc/remotes/x10/lircd.conf.mouseremote
/usr/share/lirc/remotes/x10/lircd.conf.mp3anywhere
/usr/share/lirc/remotes/ea65
/usr/share/lirc/remotes/ea65/lircd.conf.ea65
/usr/share/lirc/remotes/ea65/lircrc.ea65
/usr/share/lirc/remotes/generic
/usr/share/lirc/remotes/generic/DENON.conf
/usr/share/lirc/remotes/generic/NEC.conf
/usr/share/lirc/remotes/generic/RC-5.conf
/usr/share/lirc/remotes/generic/RCMM-32.conf
/usr/share/lirc/remotes/generic/SONY12.conf
/usr/share/lirc/remotes/generic/NEC-pulse.conf
/usr/share/lirc/remotes/generic/RECS80.conf
/usr/share/lirc/remotes/generic/MOTOROLA.conf
/usr/share/lirc/remotes/generic/XMP
/usr/share/lirc/remotes/generic/SONY20.conf
/usr/share/lirc/remotes/generic/SANYO.conf
/usr/share/lirc/remotes/generic/NEC-short-pulse.conf
/usr/share/lirc/remotes/generic/RC-6.conf
/usr/share/lirc/remotes/atiusb
/usr/share/lirc/remotes/atiusb/lircd.conf.atiusb
/usr/share/lirc/remotes/atiusb/lircd.conf.atilibusb
/usr/share/lirc/remotes/streamzap
/usr/share/lirc/remotes/streamzap/lircd.conf.streamzap
/usr/share/lirc/remotes/aureal
/usr/share/lirc/remotes/aureal/lircd.conf.atwf83
/usr/share/lirc/remotes/sigma_designs
/usr/share/lirc/remotes/sigma_designs/lircd.conf.realmagic
/usr/share/lirc/remotes/logitech
/usr/share/lirc/remotes/logitech/lircd.conf.logitech
/usr/share/lirc/remotes/logitech/lircmd.conf.logitech
/usr/share/lirc/remotes/cph03x
/usr/share/lirc/remotes/cph03x/lircd.conf.cph03x
/usr/share/lirc/remotes/acer
/usr/share/lirc/remotes/acer/lircd.conf.Aspire_6530G
/usr/share/lirc/remotes/devinput
/usr/share/lirc/remotes/devinput/lircd.conf.devinput
/usr/share/lirc/remotes/kworld
/usr/share/lirc/remotes/kworld/lircd.conf.kworld
/usr/share/lirc/remotes/technisat
/usr/share/lirc/remotes/technisat/lircd.conf.mediafocusI
/usr/share/lirc/remotes/creative
/usr/share/lirc/remotes/creative/lircd.conf.breakoutbox
/usr/share/lirc/remotes/creative/lircd.conf.alsa_usb
/usr/share/lirc/remotes/creative/lircd.conf.infracd
/usr/share/lirc/remotes/creative/lircd.conf.creative
/usr/share/lirc/remotes/creative/lircd.conf.livedrive
/usr/share/lirc/remotes/awlibusb
/usr/share/lirc/remotes/awlibusb/lircd.conf.awlibusb
/usr/share/lirc/remotes/vlsystem
/usr/share/lirc/remotes/vlsystem/lircd.conf.mplay
/usr/share/lirc/remotes/ericsson
/usr/share/lirc/remotes/ericsson/lircd.conf.bte
/usr/share/lirc/remotes/samsung
/usr/share/lirc/remotes/samsung/lircd.conf.samsung
/usr/share/lirc/remotes/kanam_accent
/usr/share/lirc/remotes/kanam_accent/lircd.conf.IR-100
/usr/share/lirc/remotes/knc_one
/usr/share/lirc/remotes/knc_one/lircd.conf.knc_one
/usr/share/lirc/remotes/hauppauge
/usr/share/lirc/remotes/hauppauge/lircmd.conf.hauppauge
/usr/share/lirc/remotes/hauppauge/lircd.conf.hauppauge
/usr/share/lirc/remotes/leadtek
/usr/share/lirc/remotes/leadtek/lircd.conf.RM-0007
/usr/share/lirc/remotes/leadtek/lircmd.conf.RM-0007
/usr/share/lirc/remotes/leadtek/lircmd.conf.RM-0010
/usr/share/lirc/remotes/leadtek/lircd.conf.RM-0010
/usr/share/lirc/remotes/leadtek/lircd.conf.PVR2000
/usr/share/lirc/remotes/chronos
/usr/share/lirc/remotes/chronos/lircd.conf.chronos
/usr/share/lirc/remotes/silitek
/usr/share/lirc/remotes/silitek/lircmd.conf.silitek
/usr/share/lirc/remotes/silitek/lircd.conf.silitek
/usr/share/lirc/remotes/hercules
/usr/share/lirc/remotes/hercules/lircd.conf.smarttv_stereo
/usr/share/lirc/remotes/pinnacle_systems
/usr/share/lirc/remotes/pinnacle_systems/lircd.conf.pctv
/usr/share/lirc/remotes/philips
/usr/share/lirc/remotes/philips/lircd.conf.srm7500libusb
/usr/share/lirc/remotes/adstech
/usr/share/lirc/remotes/adstech/lircd.conf.usbx-707
/usr/share/lirc/remotes/dvico
/usr/share/lirc/remotes/dvico/lircd.conf.fusionHDTV
/usr/share/lirc/remotes/provideo
/usr/share/lirc/remotes/provideo/lircd.conf.pv951
/usr/share/lirc/remotes/bestbuy
/usr/share/lirc/remotes/bestbuy/lircmd.conf.bestbuy2
/usr/share/lirc/remotes/bestbuy/lircmd.conf.bestbuy
/usr/share/lirc/remotes/bestbuy/lircd.conf.bestbuy2
/usr/share/lirc/remotes/bestbuy/lircd.conf.bestbuy
/usr/share/lirc/remotes/iodata
/usr/share/lirc/remotes/iodata/lircd.conf.gvbctv5pci
/usr/share/lirc/remotes/packard_bell
/usr/share/lirc/remotes/packard_bell/lircd.conf.packard_bell
/usr/share/lirc/remotes/caraca
/usr/share/lirc/remotes/caraca/lircrc.caraca
/usr/share/lirc/remotes/caraca/lircd.conf.caraca
/usr/share/lirc/remotes/README
/usr/share/lirc/remotes/tekram
/usr/share/lirc/remotes/tekram/lircd.conf.m230
/usr/share/lirc/remotes/tekram/lircmd.conf.m230
/usr/share/lirc/remotes/sasem
/usr/share/lirc/remotes/sasem/lircd.conf.sasem
/usr/share/lirc/remotes/animax
/usr/share/lirc/remotes/animax/lircd.conf.animax
/usr/share/lirc/remotes/animax/lircmd.conf.animax
/usr/share/lirc/remotes/mceusb
/usr/share/lirc/remotes/mceusb/lircd.conf.mceusb
/usr/share/lirc/remotes/adaptec
/usr/share/lirc/remotes/adaptec/lircd.conf.AVC-2410
/usr/share/lirc/remotes/apple
/usr/share/lirc/remotes/apple/lircd.conf.macmini
/usr/share/lirc/remotes/asus
/usr/share/lirc/remotes/asus/lircd.conf.asus
/usr/share/lirc/remotes/asus/lircd.conf.asusdh
/usr/share/lirc/remotes/asus/lircd.conf.digimatrix
/usr/share/lirc/remotes/pcmak
/usr/share/lirc/remotes/pcmak/lircd.conf.pcmak
/usr/share/lirc/remotes/pcmak/lircmd.conf.pcmak
/usr/share/lirc/remotes/pcmak/lircrc.pcmak
/usr/share/lirc/remotes/avermedia
/usr/share/lirc/remotes/avermedia/lircmd.conf.avermedia98
/usr/share/lirc/remotes/avermedia/lircd.conf.avermedia
/usr/share/lirc/remotes/avermedia/lircd.conf.avermedia98
/usr/share/lirc/remotes/avermedia/lircd.conf.vdomate
/usr/share/lirc/remotes/avermedia/lircmd.conf.avermedia
/usr/share/lirc/remotes/imon
/usr/share/lirc/remotes/imon/lircd.conf.imon-knob
/usr/share/lirc/remotes/imon/lircd.conf.imon-antec-veris
/usr/share/lirc/remotes/imon/lircd.conf.imon-pad
/usr/share/lirc/remotes/imon/lircd.conf.imon-rsc
/usr/share/lirc/remotes/imon/lircd.conf.imon
/usr/share/lirc/remotes/imon/lircd.conf.imon-2.4g
/usr/share/lirc/remotes/imon/lircd.conf.imon-ultrabay
/usr/share/lirc/remotes/imon/lircd.conf.imon-mceusb
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/irrecord.1.gz
/usr/share/man/man1/lircrcd.1.gz
/usr/share/man/man1/irw.1.gz
/usr/share/man/man1/irexec.1.gz
/usr/share/man/man1/irsend.1.gz
/usr/share/man/man1/irman2lirc.1.gz
/usr/share/man/man1/irpty.1.gz
/usr/share/man/man1/mode2.1.gz
/usr/share/man/man1/ircat.1.gz
/usr/share/man/man8
/usr/share/man/man8/lircd.8.gz
/usr/share/man/man8/lircmd.8.gz
/usr/share/doc
/usr/share/doc/lirc
/usr/share/doc/lirc/TODO
/usr/share/doc/lirc/copyright
/usr/share/doc/lirc/html
/usr/share/doc/lirc/html/alsa-usb.html
/usr/share/doc/lirc/html/pronto2lirc.html
/usr/share/doc/lirc/html/lircrcd.html
/usr/share/doc/lirc/html/irrecord.html
/usr/share/doc/lirc/html/imon.html
/usr/share/doc/lirc/html/configure.html
/usr/share/doc/lirc/html/imon-24g.html
/usr/share/doc/lirc/html/ircat.html
/usr/share/doc/lirc/html/irw.html
/usr/share/doc/lirc/html/audio.html
/usr/share/doc/lirc/html/irexec.html
/usr/share/doc/lirc/html/irsend.html
/usr/share/doc/lirc/html/index.html
/usr/share/doc/lirc/html/table.html
/usr/share/doc/lirc/html/smode2.html
/usr/share/doc/lirc/html/lircmd.html
/usr/share/doc/lirc/html/programs.html
/usr/share/doc/lirc/html/audio-alsa.html
/usr/share/doc/lirc/html/irxevent.html
/usr/share/doc/lirc/html/xmode2.html
/usr/share/doc/lirc/html/help.html
/usr/share/doc/lirc/html/irpty.html
/usr/share/doc/lirc/html/tira.html
/usr/share/doc/lirc/html/lircd.html
/usr/share/doc/lirc/html/devinput.html
/usr/share/doc/lirc/html/technical.html
/usr/share/doc/lirc/html/install.html
/usr/share/doc/lirc/html/mode2.html
/usr/share/doc/lirc/irxevent.keys
/usr/share/doc/lirc/NEWS.Debian.gz
/usr/share/doc/lirc/images
/usr/share/doc/lirc/images/audio_transmitter.png
/usr/share/doc/lirc/images/lirc.gif
/usr/share/doc/lirc/images/audio_wave.png
/usr/share/doc/lirc/images/atwf83.jpg
/usr/share/doc/lirc/images/diode.gif
/usr/share/doc/lirc/images/marb18.jpg
/usr/share/doc/lirc/ANNOUNCE
/usr/share/doc/lirc/changelog.Debian.gz
/usr/share/doc/lirc/NEWS.gz
/usr/share/doc/lirc/README.Debian
/usr/share/doc/lirc/README
/usr/share/doc/lirc/examples
/usr/share/doc/lirc/examples/lircmd.conf
/usr/share/doc/lirc/examples/lircd.conf
/usr/share/doc/lirc/examples/lircrc.gz
/usr/bin
/usr/bin/irman2lirc
/usr/bin/irpty
/usr/bin/irrecord
/usr/bin/mode2
/usr/bin/lircrcd
/usr/bin/ircat
/usr/bin/irexec
/usr/bin/irsend
/usr/bin/irw
/usr/sbin
/usr/sbin/lircd
/usr/sbin/lircmd
