# var/lib/dpkg/info/ramlog.conffiles
/etc/default/ramlog
/etc/init.d/ramlog
/etc/cron.daily/ramlog
/etc/logrotate.d/ramlog
/etc/ramlog.conf.upstart
