# var/lib/dpkg/info/sysfsutils.list
/.
/usr
/usr/bin
/usr/bin/systool
/usr/share
/usr/share/doc
/usr/share/doc/sysfsutils
/usr/share/doc/sysfsutils/README.Debian
/usr/share/doc/sysfsutils/changelog.gz
/usr/share/doc/sysfsutils/CREDITS
/usr/share/doc/sysfsutils/README
/usr/share/doc/sysfsutils/TODO
/usr/share/doc/sysfsutils/NEWS.gz
/usr/share/doc/sysfsutils/AUTHORS
/usr/share/doc/sysfsutils/changelog.Debian.gz
/usr/share/doc/sysfsutils/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/systool.1.gz
/etc
/etc/init.d
/etc/init.d/sysfsutils
/etc/sysfs.d
/etc/sysfs.conf
