# var/lib/dpkg/info/xdm.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/xdm
/etc/init.d
/etc/init.d/xdm
/etc/insserv.conf.d
/etc/insserv.conf.d/xdm
/etc/pam.d
/etc/pam.d/xdm
/etc/X11
/etc/X11/xdm
/etc/X11/xdm/xdm.options
/etc/X11/xdm/Xwilling
/etc/X11/xdm/Xstartup
/etc/X11/xdm/Xaccess
/etc/X11/xdm/Xsetup
/etc/X11/xdm/Xsession
/etc/X11/xdm/Xreset
/etc/X11/xdm/Xservers
/etc/X11/xdm/Xresources
/etc/X11/xdm/xdm-config
/etc/X11/app-defaults
/etc/X11/app-defaults/Chooser
/usr
/usr/bin
/usr/bin/xdm
/usr/share
/usr/share/doc
/usr/share/doc/xdm
/usr/share/doc/xdm/changelog.Debian.gz
/usr/share/doc/xdm/NEWS.Debian.gz
/usr/share/doc/xdm/README.Debian
/usr/share/doc/xdm/changelog.gz
/usr/share/doc/xdm/TODO.Debian
/usr/share/doc/xdm/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/xdm.options.5.gz
/usr/share/man/man1
/usr/share/man/man1/xdm.1.gz
/usr/share/X11
/usr/share/X11/xdm
/usr/share/X11/xdm/pixmaps
/usr/share/X11/xdm/pixmaps/xorg.xpm
/usr/share/X11/xdm/pixmaps/debianbw.xpm
/usr/share/X11/xdm/pixmaps/xorg-bw.xpm
/usr/share/X11/xdm/pixmaps/debian.xpm
/usr/lib
/usr/lib/X11
/usr/lib/X11/xdm
/usr/lib/X11/xdm/chooser
/usr/lib/X11/xdm/libXdmGreet.so
/var
/var/lib
/var/lib/xdm
