# etc/cron.daily/ramlog
#!/bin/sh
#save logs from ramdisk to harddrive:
/etc/init.d/ramlog restart
